# invizbox-luci/applications/luci-invizbox_tor/root/etc/rc.local
#!/bin/ash
#Copyright 2015 InvizBox Ltd
#
#Licensed under the InvizBox Shared License;
#you may not use this file except in compliance with the License.
#You may obtain a copy of the License at
#
#        https://www.invizbox.io/license.txt

sleep 10
/etc/init.d/tor enable
/etc/init.d/tor restart

/sbin/restore_default_passwords.sh
mv /etc/rc.local.normal /etc/rc.local
/sbin/wifi
/etc/init.d/dnsmasq restart
exit 0
